# deployment/zabbix-agent2_dropin.service
[Unit]
After=ensure-var-log-zabbix-exists.service
Requires=ensure-var-log-zabbix-exists.service
